# crontab
# Run backup daily at 2 AM
0 2 * * * /backup_librechat.sh >> /var/log/cron.log 2>&1 